# docker/supervisord/camera.conf
; ============================================================
; VisionGuard AI - Camera Capture Supervisor Configuration
; ============================================================
; Multi-process camera capture service.
; ============================================================

[supervisord]
nodaemon=true
logfile=/var/log/visionguard/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/visionguard
user=vguser

[program:camera]
command=python -m camera_capture.main
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/var/log/visionguard/camera.stdout.log
stderr_logfile=/var/log/visionguard/camera.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/app"

; Camera processes need clean shutdown
stopasgroup=true
killasgroup=true
