# docker/supervisord/ecs.conf
; ============================================================
; VisionGuard AI - ECS Supervisor Configuration
; ============================================================
; Event Classification Service - SINGLETON.
; Consumes from vg:ai:results stream.
; ============================================================

[supervisord]
nodaemon=true
logfile=/var/log/visionguard/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/visionguard
user=vguser

[program:ecs]
command=python -m event_classification.main
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/var/log/visionguard/ecs.stdout.log
stderr_logfile=/var/log/visionguard/ecs.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/app"

; ECS needs time to flush pending events on shutdown
stopasgroup=true
killasgroup=true
