# docker/supervisord/worker.conf
; ============================================================
; VisionGuard AI - AI Worker Supervisor Configuration
; ============================================================
; Stateless inference worker. Scalable.
; ============================================================

[supervisord]
nodaemon=true
logfile=/var/log/visionguard/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/visionguard
user=vguser

[program:worker]
command=python -m ai_worker.main
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/var/log/visionguard/worker.stdout.log
stderr_logfile=/var/log/visionguard/worker.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/app",ONNX_EXECUTION_PROVIDER="%(ENV_ONNX_EXECUTION_PROVIDER)s"

; Workers need longer shutdown for inference completion
stopasgroup=true
killasgroup=true
